# deploy/dagster-daemon.service
[Unit]
Description=Dagster daemon (sensors, schedules, run queue)
After=network.target docker.service dagster-code.service
Requires=docker.service dagster-code.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/match-making
EnvironmentFile=/root/match-making/.env
Environment=DAGSTER_HOME=/root/match-making
Environment=POSTGRES_HOST=127.0.0.1
Environment=ENVIRONMENT=production

ExecStart=/root/match-making/.venv/bin/dagster-daemon run

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
